# .config/systemd/user/rclone.service
[Unit]
Description=RClone Backup
Documentation=man:rclone

[Service]
Type=simple
ExecStart=%h/.dots/sh/remote